# scripts/com.cdgdg.ai-usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cdgdg.ai-usage</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"; cd /Users/cdgdg/.cdgdg-ai-usage &amp;&amp; /usr/bin/python3 scripts/collect_agent_usage.py &amp;&amp; /opt/homebrew/bin/npx vercel --prod --yes &gt;/tmp/cdgdg-ai-usage-vercel.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/cdgdg-ai-usage-launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cdgdg-ai-usage-launchd.err</string>
  <key>WorkingDirectory</key>
  <string>/Users/cdgdg/.cdgdg-ai-usage</string>
</dict>
</plist>
